# video-search.service
[Unit]
Description=Video Archive Search API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/video-search
ExecStart=/usr/bin/python3 -m uvicorn search_api:app --host 0.0.0.0 --port 8003
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
